# libs/ceilometer-expirer.crontab
#
## Unattended Installer (Semi Automated) for OpenStack (OCATA)
## Abdo Farag
## E-Mail: abdofarag85@gmail.com
## This installer forked and based on openstack-ocata-installer-centos7 v1.2.0.el7 "Siberian Lynx" released by:
## Reynaldo R. Martínez P. E-Mail: TigerLinux at Gmail dot com.
#
# Ceilometer Expirer Crontab
#

45 01 * * * root /usr/bin/ceilometer-expirer --config-file /etc/ceilometer/ceilometer.conf > /var/log/ceilometer-last-expire.log 2>&1
