# libs/keystone-flush.crontab
#
## Unattended Installer (Semi Automated) for OpenStack (OCATA)
## Abdo Farag
## E-Mail: abdofarag85@gmail.com
## This installer forked and based on openstack-ocata-installer-centos7 v1.2.0.el7 "Siberian Lynx" released by:
## Reynaldo R. Martínez P. E-Mail: TigerLinux at Gmail dot com.
#
# Keystone Cleaning Crontab
#

02 */1 * * * root /usr/local/bin/openstack-keystone-tokenflush.sh > /var/log/keystone-tokens-flush.log
